# packages/backrun/backrun-0.0.5.tar.gz/backrun-0.0.5/backrun.service
[Unit]
Description=backrun: run whitelisted scripts in background

[Service]
ExecStart=/usr/local/bin/backrun
# User=xenon

[Install]
WantedBy=multi-user.target
